# systemd/nanobot-agent@.service
[Unit]
Description=NanoBot Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
EnvironmentFile=-/etc/nanobot/nanobot.env
ExecStart=/bin/bash -lc 'repo="${NANOBOT_REPO:-%h/nanoBot_memU}"; [ -d "$repo" ] || repo="$HOME/nanoBot_memU"; venv="${NANOBOT_VENV:-%h/nanobot-venv}"; [ -x "$venv/bin/nanobot" ] || venv="$HOME/nanobot-venv"; cd "$repo" && "$venv/bin/nanobot" gateway'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
